# serviceConfig/smartdoor_restart.service
[Unit]
Description=Restart Smartdoor App

[Service]
Type=simple
ExecStart=/usr/bin/systemctl restart smartdoor.service
